# usr/lib/opkg/info/base-files-network.list
/etc/init.d/network
/etc/hotplug.d/iface/00-netstate
/sbin/ifdown
/usr/share/udhcpc/default.script
/lib/network/config.sh
/etc/hotplug.d/net/10-net
/sbin/ifup
/etc/hotplug.d/iface/10-routes
