# usr/lib/opkg/info/uboot-envtools.list
/usr/sbin/fw_printenv
/etc/init.d/uboot-envtools
